# dynamic-wallpaper.timer
[Unit]
Description=Run dynamic wallpaper every minute

[Timer]
OnBootSec=1s
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target
